# contrib/systemd/pbs_mom.service.in
[Unit]
Description=TORQUE pbs_mom daemon
Requires=network.target local-fs.target
Wants=syslog.target
After=network.target local-fs.target syslog.target

[Service]
Type=simple
User=root

LimitNOFILE=32768
LimitMEMLOCK=infinity
LimitSTACK=12582912

Environment=PBS_HOME=@PBS_HOME@
Environment=PBS_ARGS=
EnvironmentFile=-/etc/sysconfig/pbs_mom

ExecStart=@sbindir@/pbs_mom -F -d $PBS_HOME $PBS_ARGS

[Install]
WantedBy=multi-user.target
