# contrib/systemd/pbs_server.service.in
[Unit]
Description=TORQUE pbs_server daemon
Requires=network.target local-fs.target
Wants=rsyslog.target
After=trqauthd.service network.target local-fs.target rsyslog.target

[Service]
Type=simple
User=root

LimitCORE=infinity

# Let systemd guess the pid.
#
# PIDFile should not be set to PBS_HOME/server_priv/server.lock
# since it holds the pid of pbs_server on the active server (not the necessarily
# the local one) in an HA configuration.
GuessMainPID=yes

# There is also no way to test PBS_SERVERDB as was done in the sysvinit
# because the location can be changed via PBS_HOME enviroment.
# However, if you use the PBS_HOME as defined during compilation,
# you can create a unit /etc/systemd/system/pbs_server.service.d/test_serverdb.conf
# with content
#[Unit]
#ConditionPathExists=@PBS_HOME@/server_priv/serverdb

Environment=PBS_HOME=@PBS_HOME@
Environment=PBS_ARGS=
EnvironmentFile=-/etc/sysconfig/pbs_server

ExecStart=@sbindir@/pbs_server -F -d $PBS_HOME $PBS_ARGS

# Define kill parameters for shutdown.
#
# The server should be shutdown with a SIGTERM rather than using qterm
# (which will shutdown the active server in an HA configuration).
KillMode=control-group
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
