# contrib/systemd/trqauthd.service.in
[Unit]
Description=TORQUE trqauthd daemon
Requires=network.target local-fs.target
Wants=syslog.target
After=network.target local-fs.target syslog.target

[Service]
Type=simple
User=root

LimitCORE=infinity

ExecStart=@sbindir@/trqauthd -F
ExecStop=-@sbindir@/trqauthd -d

[Install]
WantedBy=multi-user.target
